# stage1/05-install-hat-detection/files/99-hat-detection.rules
SUBSYSTEM=="sound", ACTION=="add", RUN+="/usr/local/sbin/hat-detection"
